# stripminer.conf
# StripMiner - Bitcoin miner daemon

description "StripMiner Bitcoin daemon"
author "Jeremy T. Bouse <jbouse@debian.org>"

start on (started gdm and net-device-up IFACE!=lo)
stop on stopping gdm

respawn
respawn limit 10 5

env DEFAULTFILE=/etc/default/stripminer

script
  if [ -f "$DEFAULTFILE" ]; then
    . "$DEFAULTFILE"
  fi
  exec $MINER_EXEC $MINER_OPTS
end script
